# hstat/hstat.timer
[Unit]
Description=Person bash history statistics

[Timer]
OnCalendar=daily
Persistent=true
Unit=hstat.service

[Install]
WantedBy=timers.target
