# timeshift_usage/timeshift_usage.timer
[Unit]
Description=Timeshift usage

[Timer]
OnCalendar=daily
Persistent=true
Unit=timeshift_usage.service

[Install]
WantedBy=timers.target
